# get_geishu_uid_go.cfg
[program:get_geishu_uid_go]
command=/home/work/online/src/get_feishu_uid/get_feishu_uid
process_name=%(program_name)s
numprocs=1
directory=/home/work/online/src/get_feishu_uid/
environment=PATH="/usr/bin"
startretries=1
startsecs=5
autostart=yes
autorestart=true
redirect_stderr=true
stdout_logfile=/data/super_log/get_geishu_uid_go.log
